# docker/php/supervisord.conf
; Supervisorを使用してPHP-FPMとCronを並行して実行するための設定

[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:php-fpm]
command=/usr/sbin/php-fpm8.1 --nodaemonize
autostart=true
autorestart=true

[program:cron]
command=cron -f
autostart=true
autorestart=true
priority=20
